# files/etc/init.d/dnsmasq
#!/bin/sh /etc/rc.common

START=19
USE_PROCD=1

boot() {
    /etc/init.d/dnsmasq enable
    start
}