# files/etc/init.d/odhcpd
#!/bin/sh /etc/rc.common

START=00
STOP=00
USE_PROCD=1

boot() {
    /etc/init.d/odhcpd disable
    /etc/init.d/odhcpd stop
}