# files/etc/rc.local
#!/bin/sh

# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# /etc/init.d/odhcpd disable
# /etc/init.d/odhcpd stop
# /etc/init.d/dnsmasq enable
# /etc/init.d/dnsmasq start

/etc/init.d/hotspot_manager enable
/etc/init.d/hotspot_manager start

# The following is required to make sure that the gateway works after reboot
/etc/init.d/fix_ping enable
/etc/init.d/fix_ping start

exit 0
