# files/uci-defaults/80_mount_root
#!/bin/sh
[ ! -f /etc/first_login_done ] && touch /etc/first_login_pending && /usr/local/bin/first-login-setup
exit 0
