# files/uci-defaults/80_mount_root
#!/bin/sh
[ ! -f /etc/first_login_done ] && /usr/local/bin/first-login-setup
