# files/uci-defaults/99-first-login
#!/bin/sh
[ ! -f /etc/first_login_done ] && touch /etc/first_login_pending
exit 0