# files/uci-defaults/firstboot_disable_ipv6

# wget disable_ipv6.sh from blossom in future
# Found in: ~/TollGateNostrToolKit/files/root

chmod +x /root/disable_ipv6.sh
